# crontab/crontab.production
0 0 * * * /usr/scheduler/jobs/dump-sql.sh