# crontab/crontab.production
* 00 * * * /usr/scheduler/jobs/dump-sql.sh